# root/etc/udev/rules.d/90-kobodeck.rules
KERNEL=="eth*", ACTION=="add", RUN+="/usr/local/bin/kobodeck-run"
KERNEL=="wlan*", ACTION=="add", RUN+="/usr/local/bin/kobodeck-run"
