# root/etc/udev/rules.d/90-kobodeck.rules
KERNEL=="eth*", ACTION=="add", RUN+="/bin/sh -c '/usr/local/bin/kobodeck &'"
KERNEL=="wlan*", ACTION=="add", RUN+="/bin/sh -c '/usr/local/bin/kobodeck &'"
