# root/etc/udev/rules.d/90-readeckobo.rules
KERNEL=="eth*", ACTION=="add", RUN+="/usr/local/bin/readeckobo-run"
KERNEL=="wlan*", ACTION=="add", RUN+="/usr/local/bin/readeckobo-run"
