# root/etc/udev/rules.d/90-wallabako.rules
KERNEL=="eth*", RUN+="/usr/local/wallabako/wallabako-run" 
KERNEL=="wlan*", RUN+="/usr/local/wallabako/wallabako-run" 
