# root/etc/udev/rules.d/90-wallabako.rules
KERNEL=="eth*", RUN+="/usr/local/bin/wallabako-run"
KERNEL=="wlan*", RUN+="/usr/local/bin/wallabako-run"
