# root/etc/udev/rules.d/90-wallabako.rules
KERNEL=="eth*", ACTION=="add", RUN+="/usr/local/bin/wallabako-run"
KERNEL=="wlan*", ACTION=="add", RUN+="/usr/local/bin/wallabako-run"
